# deploy/wags-ai/systemd/nvidia-textfile.service
[Unit]
Description=Write NVIDIA GPU metrics to Alloy textfile collector
After=alloy.service
Wants=alloy.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/nvidia-textfile.sh
User=root
Group=root
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
