# deploy/wags-ai/systemd/ollama-textfile.timer
[Unit]
Description=Run ollama-textfile every 10 seconds
After=alloy.service

[Timer]
OnBootSec=20s
OnUnitActiveSec=10s
AccuracySec=1s
Unit=ollama-textfile.service

[Install]
WantedBy=timers.target
